# scripts/systemd/laravel-schedule-swaeduae.timer
[Unit]
Description=Run Laravel scheduler every minute (SwaedUAE)

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target
